# system-dotfiles/hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/bash /bin/udev-hdmi-switch"
