# system-dotfiles/hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/bash /bin/hdmi-switch"
