# system-dotfiles/hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/nil/bin/udev-hdmi-switch"
