# .config/nil/system-dotfiles/hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/bash /home/nil/.config/nil/hdmi-switch"
