# system/hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/dvt/bin/udev-hdmi"
